# 3dprinter-remote/rootfs/etc/services.d/ser2net/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Add-on: 3dprinter-remote
# Runs ser2net
# ==============================================================================

bashio::log.info "Starting mjpg-streamer.."

exec ser2net -d
