# 3dprinter-octoprint/rootfs/etc/services.d/proxy/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Add-on: 3dprinter-octoprint
# Stop reverse proxy
# ==============================================================================
