# 3dprinter-octoprint/rootfs/etc/services.d/proxy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Add-on: 3dprinter-octoprint
# Runs the reverse proxy.
# ==============================================================================

# Wait for OctoPrint WebUI to be available
bashio::net.wait_for 80

bashio::log.info "Starting proxy..."

exec caddy run --config /etc/caddy/Caddyfile --envfile /etc/caddy/env
