# camera-rtc/rootfs/etc/services.d/aiortc/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Add-on: aiortc
# Runs aiortc
# ==============================================================================
declare -a options


basio::log.info "Starting aiortc with arguments: '$options'"

if bashio::config.has_value 'webrtc_arguments'; then
    options+="${bashio::config 'webrtc_arguments'}"
fi

# Run aiortc
exec python3 /aiortc/examples/webcam/webcam.py "${options[@]}"
