# aiortc/rootfs/etc/services.d/aiortc/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Add-on: aiortc
# Runs aiortc
# ==============================================================================
declare -a options


bashio::log.info "Starting aiortc..."

# options+=$(cat /data/aiortc/aiortc.conf)

# # Enable debug mode
# if bashio::debug; then
#     options+=" --verbose"
# fi

# basio::log.info "Starting aiortc with arguments: '$options'"

# Run aiortc
exec python3 /aiortc/examples/webcam/webcam.py "${options[@]}"
