# 3dprinter-klipper/rootfs/etc/services.d/klipper_sim/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Add-on: 3dprinter-Klipper
# Runs Host Sim MCU
# ==============================================================================

bashio::log.info "Starting Host Sim MCU..."

## Run your program
exec "$ADDON_CONFIG_PATH"/bin/klipper_mcu_sim
