# klipper/rootfs/etc/services.d/proxy/run
#!/usr/bin/with-contenv bashio

exec caddy run -config /etc/caddy/Caddyfile