# 3dprinter-klipper/rootfs/etc/services.d/simulavr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Add-on: 3dprinter-Klipper
# Runs Host SimulAVR MCU
# ==============================================================================

bashio::log.info "Starting SimulAVR MCU..."

## Run your program
export PYTHONPATH="$SIMULAVR_PATH"/build/pysimulavr/
exec "$KLIPPER_PATH"/scripts/avrsim.py "$ADDON_CONFIG_PATH"/bin/klipper_mcu_simulavr.elf
