# 3dprinter-klipper/rootfs/etc/services.d/klippy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Add-on: 3dprinter-Klipper
# Runs Klippy
# ==============================================================================

bashio::log.info "Starting Klippy..."


if bashio::config.true 'host_mcu'; then
    while [ ! -e /tmp/klipper_host_mcu ]
    do
        bashio::log "Wait for klipper_host_mcu"
        sleep 5
    done
fi

## Run your program
exec "$ADDON_VENV_PATH"/klipper/bin/python "$ADDON_SRC_PATH"/klipper/klippy/klippy.py "$ADDON_CONFIG_PATH"/moonraker/config/printer.cfg -a /tmp/klippy_uds -l "$ADDON_CONFIG_PATH"/moonraker/logs/klippy.log
