# klipper/rootfs/etc/services.d/moonraker/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Moonraker..."

MOONRAKER_PYTHON=/data/moonraker-env/bin/python
MOONRAKER=/data/moonraker/moonraker/moonraker.py
MOONRAKER_CFG=/data/klipper_config/moonraker.conf
MOONRAKER_LOGFILE=/tmp/moonraker.log

mkdir -p /data/klipper_config

touch ${MOONRAKER_LOGFILE}
touch ${MOONRAKER_CFG}

exec ${MOONRAKER_PYTHON} ${MOONRAKER} -c ${MOONRAKER_CFG} -l ${MOONRAKER_LOGFILE}
